# misc/BootLog.service
[Unit]
Description=BootLog service

[Service]
ExecStart=/usr/bin/bash BootLog.sh

[Install]
WantedBy=multi-user.target
