# misc/NetTest.service
[Unit]
Description=NetTest service

[Service]
ExecStart=/usr/bin/bash NetTest.sh

[Install]
WantedBy=multi-user.target
