# misc/NetTest.service
[Unit]
Description=NetTest service

[Service]
ExecStart=/usr/bin/bash /etc/systemd/system/NetTest.sh

[Install]
WantedBy=multi-user.target
