# misc/sensors.service
[Unit]
Description=Sensors service

[Service]
ExecStart=/bin/bash /home/golf/sensors/run.sh

[Install]
WantedBy=multi-user.target 
